# package/deb/prerm
#!/bin/bash
set -e

# Save the list of active ndpresponder instances for later restoration
if [ "$1" = "upgrade" ] || [ "$1" = "remove" ]; then
    # Create a directory if it doesn't exist
    mkdir -p /var/lib/ndp-responder

    # Save the list of active instances
    systemctl list-units --type=service --state=active | grep ndpresponder@ | awk '{print $1}' > /var/lib/ndp-responder/active_instances.txt

    # Stop all active instances of ndpresponder@ service
    for instance in $(cat /var/lib/ndp-responder/active_instances.txt); do
        echo "Stopping $instance"
        systemctl stop "$instance" || true  # Prevent errors from halting the script
    done
    
    # Only disable services on removal, not on upgrade
    if [ "$1" = "remove" ]; then
        for instance in $(cat /var/lib/ndp-responder/active_instances.txt); do
            systemctl disable "$instance" || true
        done
    fi
fi

exit 0
